# config/systemd/user/polkit-dumb-agent.service
[Unit]
Description=Minimalistic polkit authentication agent.
PartOf=graphical-session.target

[Service]
ExecStart=/usr/bin/polkit-dumb-agent
Restart=always

[Install]
WantedBy=xsession.target
