# document/tailscale.service
[Unit]
Description = Start tailscale.

[Service]
Type = oneshot
ExecStart = /usr/bin/tailscaled

[Install]
WantedBy=multi-user.target
